# bt_receiver_ma/rootfs/etc/services.d/bt_receiver/finish
#!/usr/bin/execlineb -P
foreground { echo service stopped }
